# src/etl/crontab
* * * * * python tweets_etl.py >> /logfile.txt
